# etc/systemd/persona@.service
# persona@.service

[Unit]
Description = Persona server service

[Service]
ExecStart=/root/Persona/.build/x86_64-unknown-linux-gnu/release/Persona
