# etc/systemd/tcpproxy@.service
# tcpproxy@.service

[Unit]
Description = TCP proxy server service

[Service]
ExecStart=/root/Persona/tcpproxy/main.py
