# etc/systemd/udpproxy@.service
# udpproxy@.service

[Unit]
Description = UDP proxy server service

[Service]
ExecStart=/root/Persona/udpproxy/main.py
